# nanobot.service
[Unit]
Description=ROZ NanoBots v7 - Self-healing system daemon
After=multi-user.target network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/nanobot/nanobot.py
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=30
Nice=19
IOSchedulingClass=idle

[Install]
WantedBy=multi-user.target
