# nanobot.service
[Unit]
Description=ROZ NanoBots v5 - Self-healing system daemon
DefaultDependencies=no
After=sysinit.target network-online.target
Wants=network-online.target
Before=basic.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/nanobot/nanobot.py
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=sysinit.target
